# crontab
shell=/bin/bash path=/usr/local/sbin:/usr/local/bin:/usr/sbin:/usr/bin:/sbin:/bin:/usr/games:/usr/local/games:/snap/bin:/usr/local/bin/futuresboard
@reboot sleep 10; tmuxp load /root/apollo/session.yaml

0  */6 * * *    cd /root/apollo && python3 spot.py
15 */6 * * *    cd /root/apollo && python3 margin.py
30 */6 * * *    cd /root/apollo && python3 futures.py
45 */6 * * *    cd /root/apollo && python3 earn.py

45 8   * * *    cd /root/apollo && python3 transfer.py
45 16  * * *    sudo shutdown -r